# deploy/v5-memo-fullraw-index-build.service
[Unit]
Description=V5 Memo fullraw FTS index builder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/v5-memo
EnvironmentFile=/etc/v5-memo/env
Environment=PYTHONPATH=/opt/v5-memo/src
Environment=V5_MEMO_FULL_RAW_INDEX_MIN_FREE_GB=60
Environment=V5_MEMO_FULL_RAW_SHARDS=4
Environment=V5_MEMO_FULL_RAW_SHARD_WORKERS=4
Environment=V5_MEMO_FULL_RAW_SHARD_BATCH_FILES=16
Environment=V5_MEMO_FULL_RAW_SHARD_REMOTE=sb:researka-database/index/v5/fullraw-fts
ExecStart=/usr/bin/python3 -m v5_memo.fullraw_index build-upload-shards --commit-interval 5000
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
