# deploy/v5-memo-portfolio-prepare.timer
[Unit]
Description=Maintain the V5 Memo ready buffer every 15 minutes

[Timer]
OnCalendar=*-*-* *:00,15,30,45:00
RandomizedDelaySec=2min
Persistent=true
Unit=v5-memo-portfolio-prepare.service

[Install]
WantedBy=timers.target
